# services/cwonder@.service.in
[Unit]
Description=CWonder
After=network-online.target jfwonder.service
Requires=jfwonder.service network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=notify
PrivateTmp=true
EnvironmentFile=/usr/local/etc/wonder/cluster.conf
ExecStart=@bin_dir@/cwonder -c /usr/local/etc/wonder/cwonder_config.xml -r 2 -b
LimitRTPRIO=95
LimitRTTIME=infinity
LimitMEMLOCK=infinity
User=%i
