# services/cwonder@.service.in
[Unit]
Description=CWonder
After=network-online.target jfwonder.service
Wants=jfwonder.service network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=notify
PrivateTmp=true
EnvironmentFile=/etc/wonder/cluster.conf
ExecStart=@bin_dir@/cwonder -c /etc/wonder/cwonder_config.xml -r 2 -b
LimitRTPRIO=95
LimitRTTIME=infinity
LimitMEMLOCK=infinity
User=%i
