# services/cwonder.service.in
[Unit]
Description=CWonder

[Install]
# This is no multi-user.target for a user service, therefore we use default.target
WantedBy=default.target

[Service]
Type=notify
PrivateTmp=true
EnvironmentFile=/usr/local/etc/wonder/cluster.conf
ExecStart=@bin_dir@/cwonder -c /usr/local/etc/wonder/cwonder_config.xml -r 2 -b
LimitRTPRIO=95
LimitRTTIME=infinity
LimitMEMLOCK=infinity
